# scripts/deploy/dev/linux/ubuntu/booting/etc/rc.local
#!/bin/bash

# ------------------------------------------------------------------------------------------------------------------
# START
# ------------------------------------------------------------------------------------------------------------------


# ------------------------------------------------------------------------------------------------------------------
# START
# ------------------------------------------------------------------------------------------------------------------


# ------------------------------------------------------------------------------------------------------------------
# Base
# ------------------------------------------------------------------------------------------------------------------


# ------------------------------------------------------------------------------------------------------------------
# Booting
# ------------------------------------------------------------------------------------------------------------------


# ------------------------------------------------------------------------------------------------------------------
# Network
# ------------------------------------------------------------------------------------------------------------------

# >>>> Firewall

# >>>> IPv4

# >>>> IPv6

# ------------------------------------------------------------------------------------------------------------------
# Packages
# ------------------------------------------------------------------------------------------------------------------


# ------------------------------------------------------------------------------------------------------------------
# Process
# ------------------------------------------------------------------------------------------------------------------


# ------------------------------------------------------------------------------------------------------------------
# Security
# ------------------------------------------------------------------------------------------------------------------

# >>>> Directories

# >>>> Files

# >>>> Groups

# >>>> Users

# ------------------------------------------------------------------------------------------------------------------
# Utility
# ------------------------------------------------------------------------------------------------------------------

# >>>> Power

# ------------------------------------------------------------------------------------------------------------------
# END
# ------------------------------------------------------------------------------------------------------------------

exit 0
