# scripts/containers/dev/tools/supervisor/etc/supervisord.conf
[supervisord]
user=root
nodaemon=true
pidfile=/var/run/supervisor/supervisord.pid
logfile=/var/log/supervisor/supervisord.log
childlogdir=/var/log/supervisor
logfile_maxbytes=0
loglevel=info

[program:php-fpm]
command=php-fpm -F
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/nginx.out.log
stdout_logfile_maxbytes=0
stderr_logfile=/var/log/supervisor/nginx.err.log
stderr_logfile_maxbytes=0

[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/php-fpm.out.log
stdout_logfile_maxbytes=0
stderr_logfile=/var/log/supervisor/php-fpm.err.log
stderr_logfile_maxbytes=0
