# scripts/linux/ubuntu/tools/settings/rc.local
#!/bin/bash
# ======================================================================================================================
# Scripts - Linux - Ubuntu - Boot
# ======================================================================================================================
# >>>> Platform
if [ "${PLATFORM_TYPE}" == "Linux" ]; then

  # --------------------------------------------------------------------------------------------------------------------
  # Network
  # --------------------------------------------------------------------------------------------------------------------


  # --------------------------------------------------------------------------------------------------------------------
  # Firewall
  # --------------------------------------------------------------------------------------------------------------------
  echo "Firewall"

  # --------------------------------------------------------------------------------------------------------------------
  # Packages
  # --------------------------------------------------------------------------------------------------------------------


  # --------------------------------------------------------------------------------------------------------------------
  # Services
  # --------------------------------------------------------------------------------------------------------------------

fi

exit 0
