# scripts/linux/ubuntu/app/crontab/_install.sh
#!/bin/bash
# ======================================================================================================================
# Scripts - Linux - Ubuntu - Cron - Install
# ======================================================================================================================
# >>>> Platform
if [ "${PLATFORM_TYPE}" == "Linux" ]; then

  echo ">>>> Linux - Cron"
  echo

  # --------------------------------------------------------------------------------------------------------------------
  # Cron - Install the packages
  # --------------------------------------------------------------------------------------------------------------------

  # >>>> Crontab
  echo ">>>> Crontab"
  echo

  local CRONTAB_STATUS
  CRONTAB_STATUS=$(systemctl is-active cron)
  if [ "${CRONTAB_STATUS}" == "active" ]; then
    sudo systemctl stop cron.service
    echo
  fi

  # --------------------------------------------------------------------------------------------------------------------
  # Cron - Update the configuration
  # --------------------------------------------------------------------------------------------------------------------

  # >>>> crontab


  # --------------------------------------------------------------------------------------------------------------------
  # Cron - Check the configuration
  # --------------------------------------------------------------------------------------------------------------------

  CRONTAB_STATUS=$(systemctl is-active cron)
  if [ "${CRONTAB_STATUS}" == "inactive" ]; then
    sudo systemctl start cron.service
  fi
  echo

  sudo systemctl status cron.service --no-pager

fi
